# tools/supervisor.ini
[supervisord]
logfile=/dev/null
loglevel=info
nodaemon=true

[program:api]
command=node tools/tokensrv.ts

[program:react]
command=yarn start

[program:swa]
command=swa start http://localhost:3000 --api-location http://localhost:8765 --host 0.0.0.0 -s
